# crontab
SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin:/bin

# cron job to run the script every day at 4:00 AM
15 20 * * * cd /app && python today.py >> /var/log/cron.log 2>&1
